# leader_election_server.service
[Unit]
Description=Leader Election Server
After=network.target

[Service]
Type=simple
WorkingDirectory=/usr/local/bin
ExecStart=/usr/bin/python3 /usr/local/bin/leader_election_server.py
Restart=always

[Install]
WantedBy=multi-user.target
